# etc/systemd/system/media-data.mount
[Unit]
Description=xHD mount script

[Mount]
What=/dev/sda2
Where=/media/data
Type=ntfs
Options=defaults,gid=users,dmask=002,fmask=002,noauto,nofail
