# etc/systemd/system/mnt-share.mount
[Unit]
Description=xHD mount script

[Mount]
What=/dev/sda2
Where=/mnt/share
Type=ntfs-3g
Options=gid=users,dmask=002,fmask=002,noauto,nofail

